# time-machine/services/netatalk/run
#!/bin/bash

rm -f /var/lock/netatalk
netatalk -d
